# image/sd-image/seed/fs/etc/init/ttyAMA0.conf
start on stopped rc or RUNLEVEL=[12345]
stop on runlevel [!12345]

respawn
exec /sbin/getty -L 115200 ttyAMA0 vt102
